# deploy/supervisor-ipop360-ssr.conf
# Supervisor program for the Inertia SSR server (spec-063).
#
# Renders {{DEPLOY_PATH}} at deploy time (see .github/workflows/deploy.yml
# "Register the Inertia SSR supervisor program") then `supervisorctl reread +
# update` loads it. The Post-deploy step restarts it each deploy so the freshly
# built bootstrap/ssr/ssr.js is loaded (supervisor only re-spawns on config
# change, not file-content change).
#
# Fail-safe: if this process is down/crashing, Inertia's HttpGateway catches the
# connection failure and renders client-side — the site never goes down because
# of SSR. Kill-switch: INERTIA_SSR_ENABLED=false in .env + config:cache.
[program:ipop360-ssr]
process_name=%(program_name)s_%(process_num)02d
command=node {{DEPLOY_PATH}}/bootstrap/ssr/ssr.js
directory={{DEPLOY_PATH}}
autostart=true
autorestart=true
startsecs=3
startretries=8
stopwaitsecs=10
user=www-data
environment=NODE_ENV="production",INERTIA_SSR_PORT="13714"
stdout_logfile={{DEPLOY_PATH}}/storage/logs/ssr.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile={{DEPLOY_PATH}}/storage/logs/ssr.err.log
stderr_logfile_maxbytes=5MB
redirect_stderr=false
